# src/.services/sample/fsroot/etc/sv/sample/finish
# This Script will called when runit finished this Service
