# src/.services/sample/fsroot/etc/sv/sample/run
# This Script will called when runit started this Service
